# rc.local
#!/bin/bash

APP="$(dirname "$(readlink -m "$BASH_SOURCE")")"

source $HOME/.rvm/scripts/rvm
rvm use ruby-head
$APP/status &
while true
do
	stdbuf -oL -eL ruby $APP/enslaver.rb
	sleep 10
done 2>&1 | logger -t enslaver
